# data/supervisord.conf
[supervisord]
nodaemon=true

[program:chain]
command=/usr/bin/python    /opt/chain/manage.py   runserver   0.0.0.0:8001
autorestart=True
autostart=True
redirect_stderr=True

[program:webssh]
command=/usr/bin/python    /opt/chain/webssh/main.py
autorestart=True
autostart=True
redirect_stderr=True


[program:celery]
command=/usr/bin/python    /opt/chain/manage.py  celery worker  -c  4        --loglevel=info
autorestart=True
autostart=True
redirect_stderr=True